# scripts/systemd/yubing-admin.service
[Unit]
Description=YuBing Insight Vue Admin dev server (:5173)
After=yubing-api.service network-online.target
Wants=yubing-api.service network-online.target

[Service]
Type=simple
User=root
WorkingDirectory=/root/workspace/YuBingInsight
Environment=NODE_OPTIONS=--max-old-space-size=2048
Environment=PATH=/root/.nvm/versions/node/v20.20.2/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/root/.nvm/versions/node/v20.20.2/bin/pnpm dev:admin
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=yubing-admin

[Install]
WantedBy=multi-user.target
